# profile/airootfs/etc/systemd/system/automount.service
[Unit]
Description=Automount disks

[Service]
Type=oneshot
ExecStart=/usr/local/bin/automount
#ExecStop=

[Install]
WantedBy=multi-user.target
